# deployment/mexc-screener.service
[Unit]
Description=Screener: volume filter + SWH/SWL sweep (Linux service, runs every 10 min)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ultimate-screener
EnvironmentFile=/opt/ultimate-screener/.env
ExecStart=/opt/ultimate-screener/.venv/bin/python -u src/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
